# dot_config/systemd/user/anki.service
[Unit]
Description=Anik

[Service]
ExecStart=/usr/bin/anki
TimeoutSec=10

[Install]
WantedBy=xsession.target