# dot_config/systemd/user/arandr.service
[Unit]
Description=Load Arandr config

[Service]
ExecStart=%h/.screenlayout/default.sh
Type=oneshot

[Install]
WantedBy=xsession.target

