# dot_config/systemd/user/copyq.service
[Unit]
Description=Copyq

[Service]
ExecStart=/usr/bin/copyq

[Install]
WantedBy=xsession.target