# dot_config/systemd/user/emacs-misc.service
[Unit]
Description=Emacs Misc

[Service]
Type=forking
ExecStart=/usr/bin/emacs --daemon=misc

[Install]
WantedBy=xsession.target