# dot_config/systemd/user/emacs-misc.service
[Unit]
Description=Emacs Misc

[Service]
ExecStart=/usr/bin/emacs --fg-daemon=misc

[Install]
WantedBy=xsession.target