# dot_config/systemd/user/fcitx5.service
[Unit]
Description=Fcitx5

[Service]
ExecStart=/usr/bin/fcitx5

[Install]
WantedBy=xsession.target