# dot_config/systemd/user/goldendict.service
[Unit]
Description=Goldendict

[Service]
Environment="QT_SCREEN_SCALE_FACTORS=1"
ExecStart=/usr/bin/goldendict

[Install]
WantedBy=xsession.target