# dot_config/systemd/user/goldendict.service
[Unit]
Description=Goldendict

[Service]
ExecStart=/usr/bin/goldendict

[Install]
WantedBy=xsession.target