# dot_config/systemd/user/i3-session.target
[Unit]
Description=A hook target where we know i3 have been started

[Install]
WantedBy=suspend.target hibernate.target