# dot_config/systemd/user/imwheel.service
[Unit]
Description=IMWheel

[Service]
ExecStart=/usr/bin/imwheel -d -b '4567'
ExecStop=/usr/bin/kill imwheel

[Install]
WantedBy=xsession.target