# dot_config/systemd/user/keepassxc.service
[Unit]
Description=Keepassxc

[Service]
ExecStart=/usr/bin/keepassxc

[Install]
WantedBy=xsession.target