# dot_config/systemd/user/nm-applet.service
[Unit]
Description=NM Applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=xsession.target