# dot_config/systemd/user/openrgb.service.tmpl
[Unit]
Description=OpenRGB

[Service]
Type=oneshot
ExecStart=/usr/bin/openrgb -p nolight

{{if (eq .chezmoi.hostname "kazuki")}}

[Install]
WantedBy=xsession.target
{{end}}