# dot_config/systemd/user/picom.service
[Unit]
Description=Compositing manager

[Service]
ExecStart=/usr/bin/picom
Restart=on-failure
RestartSec=30

[Install]
WantedBy=xsession.target
