# dot_config/systemd/user/playerctld.service
[Unit]
Description=Playerctld

[Service]
ExecStart=/usr/bin/playerctld

[Install]
WantedBy=xsession.target