# dot_config/systemd/user/redshift-gtk.service
[Unit]
Description=Redshift-Gtk

[Service]
ExecStart=/usr/bin/redshift-gtk

[Install]
WantedBy=xsession.target