# dot_config/systemd/user/syncthing-gtk.service
[Unit]
Description=syncthing-Gtk

[Service]
ExecStart=/usr/bin/syncthing-gtk

[Install]
WantedBy=xsession.target