# dot_config/systemd/user/wallpaper.service
[Unit]
Description=Set Wallpaper

[Service]
ExecStart=/usr/bin/nitrogen --restore
Type=oneshot

[Install]
WantedBy=xsession.target

