# dot_config/systemd/user/xsession.target
[Unit]
Description=Sevices which requires X11
Requires=graphical.target
After=graphical.target suspend.target hibernate.target
