# dot_config/systemd/user/xsession.target
[Unit]
Description=Sevices which requires X11
Requisite=i3-session.target
After=i3-session.target

[Install]
RequiredBy=i3-session.target