# dot_config/systemd/user/docker.service.d/override.conf.tmpl
[Service]
Environment="DOCKERD_ROOTLESS_ROOTLESSKIT_PORT_DRIVER=slirp4netns"
